# files/firstboot-ssh.service
[Unit]
Description=One-time SSH bootstrap (first boot only)
ConditionFirstBoot=yes
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/firstboot-ssh.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

